# usr/lib/systemd/system/temp_monitor.service
[Unit]
Description=PC Temperature Monitor
After=network.target
[Service]
ExecStart=/usr/local/bin/temp_monitor_service.sh
Restart=always
User=root
Group=root
[Install]
WantedBy=multi-user.target